# docker/pimcore/s6/etc/services.d/php-fpm/run
#!/usr/bin/with-contenv sh

mkdir -p /run/php
exec /usr/sbin/php-fpm8.1 --nodaemonize --fpm-config /etc/php/8.1/fpm/php-fpm.conf
