# docker/pimcore/s6/etc/services.d/php-fpm/run
#!/usr/bin/with-contenv sh

mkdir -p /run/php
exec /usr/sbin/php-fpm7.3 --nodaemonize --fpm-config /etc/php/7.3/fpm/php-fpm.conf
